# Sources/Extras/CoreMediaIO/DeviceAbstractionLayer/Devices/Sample/Assistant/com.apple.cmio.DPA.Sample.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.cmio.DPA.Sample</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/CoreMediaIO/Plug-Ins/DAL/Sample.plugin/Contents/Resources/SampleAssistant</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cmio.DPA.Sample</key>
		<true/>
	</dict>
</dict>
</plist>
